# systemd/industrial-scanner-logger-api.service
[Unit]
Description=Industrial Scanner Logger REST API
Documentation=file:@INSTALL_DIR@/README.md
Wants=network-online.target
After=network-online.target postgresql.service

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_DIR@
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PYTHONUNBUFFERED=1
ExecStart=@PYTHON_BIN@ -m industrial_scanner_logger.api
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=10
NoNewPrivileges=true
UMask=0027
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
